# crontab
*/5 * * * * python /wwwroot/wepython/timer_send.py